# AI/shared_hub/systemd-user/projeler-ai-hub.service
[Unit]
Description=Projeler AI Hub Web
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/Projeler/.ai-hub
ExecStart=%h/Projeler/.ai-hub/start_projeler_ai_hub.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
